# dit4c-gatehouse/src/main/docker/etc/supervisord.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
redirect_stderr=true